# conf/fsync.ini
[program:fsync]
command=sh /var/gunicorn/runinenv /var/gunicorn/venvs/fsync /var/gunicorn/venvs/fsync/bin/gunicorn -c /var/gunicorn/apps/fsync/conf/gunicorn..conf.py djprojcet.wsgi:application
directory=/var/gunicorn/apps/fsync/djproject
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/supervisor/fsync_stdout.log
stderr_logfile=/var/log/supervisor/fsync_stderr.log
