# conf/prod/supervisor/fsync.ini

[program:fsync]
command=sh /var/gunicorn/runinenv /var/gunicorn/venvs/fsync /var/gunicorn/venvs/fsync/bin/gunicorn -c /var/gunicorn/apps/fsync/conf/prod/gunicorn.conf.py djproject.wsgi:application
directory=/var/gunicorn/apps/fsync/djproject
autostart=true
autorestart=unexpected
redirect_stderr=true
stdout_logfile=/var/log/supervisor/fsync_stdout.log
stderr_logfile=/var/log/supervisor/fsync_stderr.log

